# scripts/pi/worfbot-gateway.service
[Unit]
Description=Worfbot Gateway
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/actions-runner/_work/zeal/zeal/apps/worfbot-gateway
ExecStart=/usr/bin/node dist/main.js
EnvironmentFile=/etc/worfbot-gateway/env
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
